# systemd/litestream.service
[Unit]
Description=Litestream SQLite replication
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=app
Group=app
EnvironmentFile=/etc/litestream/env
ExecStart=/usr/local/bin/litestream replicate -config /etc/litestream.yml
Restart=always
RestartSec=5

NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/app /var/log/app

[Install]
WantedBy=multi-user.target
